# roles/ubuntu/files/ddns-update.service
[Unit]
Description=Update IEServer DDNS

[Service]
Type=simple
ExecStart=/opt/local/ddns-update/update.sh

[Install]
WantedBy=default.target
