# roles/ubuntu/files/gpg-agent.service
[Unit]
Description=GnuPG agent

[Service]
Type=simple
ExecStart=/usr/bin/gpg-agent

[Install]
WantedBy=default.target
