# recipes/ubuntu/files/home/k0kubun/.config/autostart/nm-applet.desktop
[Desktop Entry]
Name=Network Manager Applet
Type=Application
Exec=nm-applet
Terminal=false
