# recipes/ubuntu/files/ddns-update.service
[Unit]
Description=Update Route53 DDNS

[Service]
Type=simple
ExecStart=/usr/bin/ruby /opt/local/ddns-update/update

[Install]
WantedBy=default.target
