# recipes/ubuntu/files/ddns-update.timer
[Unit]
Description=Timer to update Route53

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
Unit=ddns-update.service

[Install]
WantedBy=timers.target
